# scheduling/com.liquorbond.automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liquorbond.automation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_ME_PATH_TO_SCRIPT/smart_dispatcher.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <dict>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Minute</key>
            <integer>45</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>REPLACE_ME_PATH_TO_PROJECT/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_ME_PATH_TO_PROJECT/launchd_stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
